# packaging/systemd/zensight-sensor-netlink.service
[Unit]
Description=ZenSight netlink — Linux kernel net telemetry (RTNETLINK/sock_diag)
Documentation=https://github.com/p13marc/zensight
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=/usr/bin/zensight-sensor-netlink --config /etc/zensight/netlink.json5
Restart=on-failure
RestartSec=5s
# Graceful stop (#161): SIGTERM lets the process publish its offline status and
# tombstone any firing alerts before exit.
KillSignal=SIGTERM
TimeoutStopSec=20s
# Run under a transient system user with a minimal sandbox.
DynamicUser=yes
NoNewPrivileges=yes
ProtectSystem=strict
ProtectHome=yes
ReadOnlyPaths=/etc/zensight

# DEMO POSTURE — "just run" with every collector enabled.
#
# The default reads (interfaces/routes/neighbors/addresses/sockets/ethtool/tc/
# diagnostics/RTNETLINK events/XFRM SA dump) are unprivileged and need NONE of
# the below. These capabilities unlock the *optional* collectors so a demo shows
# everything:
#   - CAP_NET_ADMIN → nftables + conntrack telemetry, and subscribing to the XFRM
#     monitor multicast groups (real-time IPsec lifecycle → events/ipsec/*).
#   - CAP_BPF (+ CAP_PERFMON for perf-based attach) → the eBPF module (#114).
#
# eBPF additionally requires a binary built with `--features ebpf` and
# `collect.ebpf = true`; the stock release binary ignores the BPF caps.
# (Ambient caps coexist with NoNewPrivileges=yes — systemd raises them at exec.)
# To return to the locked-down baseline, drop the three lines below.
AmbientCapabilities=CAP_NET_ADMIN CAP_BPF CAP_PERFMON
CapabilityBoundingSet=CAP_NET_ADMIN CAP_BPF CAP_PERFMON
ProtectKernelTunables=no

[Install]
WantedBy=multi-user.target
